# private_Library/LaunchAgents/com.ameba.swiftbar.autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Launch SwiftBar at login, reproducibly. SwiftBar's own "Launch at Login"
     toggle uses SMAppService (no defaults key, GUI-only), so this LaunchAgent
     is the chezmoi-trackable equivalent. RunAtLoad fires once at login (and on
     `launchctl bootstrap`); no KeepAlive — SwiftBar manages its own lifetime.
     `open -a` resolves the app via LaunchServices, so the path isn't hardcoded.
     Loaded for the current session by run_onchange_after_setup-swiftbar. -->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ameba.swiftbar.autostart</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/open</string>
		<string>-a</string>
		<string>SwiftBar</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/dev/null</string>
	<key>StandardErrorPath</key>
	<string>/dev/null</string>
</dict>
</plist>
